# systemd/minecraft-map@.timer
[Unit]
Description=Minecraft map render timer for %i

[Timer]
OnCalendar=daily
Unit=minecraft-map@%i.service

[Install]
WantedBy=default.target
