# units/minecraft-map@.service
[Unit]
Description=Minecraft map render of %I

[Service]
Type=simple
WorkingDirectory=%h/%i
ExecStart=%h/scripts/map.sh
SyslogIdentifier=render-%i
