# units/minecraft-backup@.timer
[Unit]
Description=Minecraft backup timer for %i

[Timer]
OnCalendar=daily
Unit=minecraft-backup@%i.service

[Install]
WantedBy=timers.target
